# deploy/game125.service
[Unit]
Description=Game 125 (WebSocket server)
After=network.target

[Service]
WorkingDirectory=%h/BRIDGE
ExecStart=/usr/bin/node server/index.js
Environment=PORT=3000
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
